# 99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="$PATH_TO_SCRIPT/hotplug.sh"
