# 99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/gehidore/.local/pbin/udev-monitor-hotplug/hotplug.sh"
